# supervisor/conf.d/apache2.conf
[program:apache2]
command=apachectl -D "FOREGROUND" -k start
killasgroup=true
stopasgroup=true
autorestart=true
redirect_stderr=true