# supervisor/conf.d/mysql.conf
[program:mysql]
command=/entrypoint.sh mysqld 
killasgroup=true
stopasgroup=true
autorestart=true